# services/uwsgi/altrepo-api.socket
[Unit]
Description=uWSGI socket for ALTRepo API

[Socket]
ListenStream=/run/altrepo-api.sock

[Install]
WantedBy=sockets.target
